# root/etc/s6-overlay/s6-rc.d/medusa-cache-refresh/up
/etc/s6-overlay/s6-rc.d/medusa-cache-refresh/run
